# php5.6-ubuntu16.04/.docker/php5.6-apache-supervisor/crontab
*/1 * * * * root php /var/www/html/bin/console crontasks:run >> /tmp/info
